# systemd/system/web-apps.target
[Unit]
After=network.target
Wants=web-apps-http.target web-apps-worker.target

[Install]
WantedBy=multi-user.target
